# templates/etc/systemd/system/ctrl-alt-del-capture.service.erb
[Unit]
Description=Log calls to Ctrl-Alt-Delete - Managed by Puppet

[Service]
Type=simple
StandardOutput=syslog
SyslogIdentifier=ctrl-alt-del
SyslogFacility=<%= @facility %>
SyslogLevel=<%= @severity %>
ExecStart=<%= @_log_cmnd %>

[Install]
WantedBy=ctrl-alt-del.target
